# Docker/Matsya/1.0/vbox-client
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="VirtualBox Guest GUI Integration Service"

depend() {
        need virtualbox-guest-additions
}

start() {
        ebegin "Starting VirtualBox Guest GUI Integration Service"
        /usr/sbin/VBoxClient --draganddrop
        /usr/sbin/VBoxClient --vmsvga
        eend $?
}

stop() {
        ebegin "Stopping VirtualBox Guest GUI Integration Service"
        killall VBoxClient
        eend $?
}
